# scripts/batch/systemd/bizi-announcement-update.service
[Unit]
Description=Bizi 지원사업 공고 일일 배치 갱신
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
User=ec2-user
Group=docker
WorkingDirectory=/home/ec2-user/bizi
ExecStart=/usr/bin/docker compose -f docker-compose.prod.yaml run --rm batch-updater
TimeoutStartSec=7200
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
