# scripts/batch/systemd/bizi-announcement-update.timer
[Unit]
Description=Bizi 공고 갱신 일일 타이머

[Timer]
OnCalendar=*-*-* 18:30:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
